# deployment/docker/h2o/templates/supervisord.conf
[supervisord]
nodaemon=true

[program:h2o]
user=h2o
startsecs = 0
autorestart = false
command= java -Xmx500m -jar /opt/h2o/h2o.jar -port 33333 -log_dir /var/log/h2o