# deployment/docker/openrefine/templates/supervisord.conf
[supervisord]
nodaemon=true

[program:openrefine]
user=openrefine
startsecs = 0
autorestart = false
command= /opt/openrefine/refine -p PORT -i HOST