# deployment/docker/nifi/templates/supervisord.conf
[supervisord]
nodaemon=true

[program:nifi]
user=nifi
startsecs = 0
autorestart = false
command= /opt/nifi/bin/nifi.sh run